# deploy/systemd/v2board-worker.service
[Unit]
Description=V2Board background worker
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=60s
StartLimitBurst=5

[Service]
Type=notify
NotifyAccess=main
WatchdogSec=30s
User=v2board-worker
Group=v2board-worker
WorkingDirectory=/var/lib/v2board/worker

Environment=V2BOARD_ENV=production
Environment=V2BOARD_RUNTIME_ROOT=/var/lib/v2board
Environment=V2BOARD_CONFIG_PATH=/var/lib/v2board/worker/config.json
Environment=V2BOARD_WORKER_HEALTH_FILE=/run/v2board-worker/health
Environment=V2BOARD_WORKER_HEARTBEAT_INTERVAL_SECONDS=10
Environment=V2BOARD_WORKER_SHUTDOWN_TIMEOUT_SECONDS=30
Environment=RUST_LOG=v2board_workers=info

ExecStart=/opt/v2board/current/bin/v2board-workers
Restart=on-failure
RestartSec=3s
TimeoutStartSec=60s
TimeoutStopSec=45s
KillSignal=SIGTERM
UMask=0077

RuntimeDirectory=v2board-worker
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=no

StandardOutput=journal
StandardError=journal
SyslogIdentifier=v2board-worker

NoNewPrivileges=true
CapabilityBoundingSet=
AmbientCapabilities=
PrivateDevices=true
PrivateMounts=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProtectSystem=strict
ProcSubset=pid
LockPersonality=true
MemoryDenyWriteExecute=true
RemoveIPC=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
DevicePolicy=closed
ReadOnlyPaths=/opt/v2board/current
ReadWritePaths=/var/lib/v2board/worker
LimitNOFILE=65536
TasksMax=4096
# Background jobs yield to the API under contention (CPUWeight below the API's
# 100); memory ceilings bound a runaway batch without a hard CPU quota.
MemoryHigh=512M
MemoryMax=1G
CPUWeight=50

[Install]
WantedBy=multi-user.target
