# deploy/systemd/v2board-api.service
[Unit]
Description=V2Board API and frontend service
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=60s
StartLimitBurst=5

[Service]
Type=exec
User=v2board-api
Group=v2board-api
WorkingDirectory=/var/lib/v2board/api

Environment=V2BOARD_ENV=production
Environment=V2BOARD_RUNTIME_ROOT=/var/lib/v2board
Environment=V2BOARD_CONFIG_PATH=/var/lib/v2board/api/config.json
Environment=V2BOARD_RULE_DIR=/var/lib/v2board/rules
Environment=V2BOARD_FRONTEND_DIR=/opt/v2board/current/frontend
Environment=RUST_LOG=v2board_api=info,tower_http=info

ExecStart=/opt/v2board/current/bin/v2board-api
Restart=on-failure
RestartSec=3s
TimeoutStartSec=60s
TimeoutStopSec=60s
KillSignal=SIGTERM
UMask=0077

StandardOutput=journal
StandardError=journal
SyslogIdentifier=v2board-api

NoNewPrivileges=true
CapabilityBoundingSet=
AmbientCapabilities=
PrivateDevices=true
PrivateMounts=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProtectSystem=strict
ProcSubset=pid
LockPersonality=true
MemoryDenyWriteExecute=true
RemoveIPC=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
DevicePolicy=closed
ReadOnlyPaths=/opt/v2board/current /var/lib/v2board/rules
ReadWritePaths=/var/lib/v2board/api
LimitNOFILE=65536
TasksMax=4096

[Install]
WantedBy=multi-user.target
